# services/linux-systemd.service
[Unit]
Description=Star-DART GitHub Star Poller
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/youruser/star-dart
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /home/youruser/star-dart/scripts/webhook_poller.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
